# scripts/launchd/com.rich.jin10-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rich.jin10-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/rich/jin10-monitor/scripts/run_monitor.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/rich/jin10-monitor</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/rich/jin10-monitor/logs/jin10-monitor.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/rich/jin10-monitor/logs/jin10-monitor.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
